# lib/linux/template.upstart
description  "{{name}}"
author       "{{author}}"

start on startup
stop on runlevel [016]
kill signal QUIT

{{#wait}}
post-stop exec sleep {{wait}}
{{/wait}}

respawn

{{#user}}
setuid {{user}}
{{/user}}

{{#env}}
env {{key}}={{value}}
{{/env}}

script
  chdir {{path}}
  exec {{bin}}
end script
