# lib/linux/template.upstart
description  "{{name}}"
author       "{{author}}"

start on startup
stop on shutdown
kill signal QUIT

respawn
respawn limit 20 5

{{#user}}
setuid {{user}}
{{/user}}

{{#env}}
env {{key}}={{value}}
{{/env}}

script
  chdir {{path}}
  exec {{bin}}
end script
